# deploy/systemd/frogswork-api.service
[Unit]
Description=FrogsWork File Storage API
After=network-online.target
Wants=network-online.target
RequiresMountsFor=/data

[Service]
Type=simple
User=frogswork
Group=frogswork
WorkingDirectory=/opt/frogswork/backend
Environment=PYTHONUNBUFFERED=1
ExecStart=/opt/frogswork/venv/bin/uvicorn frogswork_api.main:app --host 127.0.0.1 --port 8000
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
